# crontab
0 1 * * * flock -n /home/greg/greg.lock /usr/bin/greg sync 
