# crontab
0 * * * * flock -n /home/greg/greg.lock /usr/bin/greg sync 
